# network/plugins/tr181-dnssd/files/etc/uci-defaults/99-set-lla-controller-dnssd
#!/bin/sh

sed -i 's|Controller = "mod-umdns-uci"|Controller = "mod-umdns"|' /etc/amx/tr181-dnssd/tr181-dnssd_definition.odl

exit 0
