# network/plugins/time-manager/files/etc/uci-defaults/99-disable-chronyd
#!/bin/sh

/etc/init.d/chronyd disable

exit 0
