# network/utilities/tr181-rest-api/files/etc/uci-defaults/99-change-default-webserver-port
#!/bin/sh

sed -i 's|Port = 80;|Port = 8080;|' /etc/amx/tr181-httpaccess/defaults.d/00_lan_access_ui.odl

exit 0
