# network/plugins/tr181-httpaccess/files/etc/uci-defaults/99-change-lighttpd-user-and-group
#!/bin/sh

# fix amx-fcgi fcgi socket handling
sed -i 's|server.username             = "http"|server.username             = "lighttpd"|' /etc/lighttpd/lighttpd.template
sed -i 's|server.groupname            = "www-data"|server.groupname            = "lighttpd"|' /etc/lighttpd/lighttpd.template

exit 0
