# network/plugins/dhcpv4-manager/files/etc/uci-defaults/99-set-lla-controller
#!/bin/sh

sed -i 's|Controller = "mod-dhcp-uci"|Controller = "mod-dhcp-dnsmasq"|' /etc/amx/dhcpv4-manager/server/dhcpv4-manager-server.odl

exit 0
